# fix-Logitech-HD-Pro-Webcam-C920.service
[Service]
ExecStart=/usr/local/bin/fix-Logitech-HD-Pro-Webcam-C920.sh
Type=simple

[Install]
WantedBy=default.target
